# files/dashboard/tailpipe-dashboard.service
[Unit]
Description=TailPiPE LCD Dashboard
After=network-online.target NetworkManager.service tailscaled.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tailpipe-dashboard
# Auto-detects the LCD fbdev; override with TAILPIPE_FB=/dev/fbN if needed.
Environment=TAILPIPE_FB=auto
# If no LCD is found, the script exits 0 and systemd stops restarting.
Restart=on-failure
RestartSec=5
User=root
# Ensure evdev and framebuffer access
SupplementaryGroups=input video

[Install]
WantedBy=multi-user.target
