# files/ts-hosts-update.service
[Unit]
Description=Refresh /etc/hosts.tailscale from tailscale status
After=tailscaled.service
Requires=tailscaled.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ts-hosts-update
